# roles/supervisor/templates/conf.d/programs.conf.j2
{% for program in hostprogram|sort %}
[program:{{program}}]
{% for key in hostprogram[program]|sort %}
{{ key }} = {{ hostprogram[program][key] }}
{% endfor %}

{% endfor %}